# roles/pipewire_audio/templates/mic-dsp-link.service.j2
[Unit]
Description=Link hardware microphone to RNNoise DSP filter
After=wireplumber.service
Wants=wireplumber.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c '\
for i in 1 2 3 4 5; do \
  HW_MIC=$$(pactl list short sources 2>/dev/null | grep alsa_input | grep -v monitor | head -1 | awk "{print $$2}"); \
  [ -n "$$HW_MIC" ] && break; \
  sleep 1; \
done; \
[ -n "$$HW_MIC" ] || exit 0; \
pactl set-source-volume $$HW_MIC {{ mic_hardware_volume_pct | default(100) }}% 2>/dev/null || true; \
pw-link $$HW_MIC:capture_FL effect_input.rnnoise:input_FL 2>/dev/null; \
pw-link $$HW_MIC:capture_FR effect_input.rnnoise:input_FR 2>/dev/null; \
pactl set-default-source effect_output.rnnoise 2>/dev/null || true; \
exit 0'
RemainAfterExit=yes

[Install]
WantedBy=default.target
