# dot_config/systemd/user/chezmoi-update.timer
[Unit]
Description=Run chezmoi dotfiles update daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
