# dot_config/systemd/user/chezmoi-update.service
[Unit]
Description=chezmoi dotfiles auto-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/linuxbrew/.linuxbrew/bin/chezmoi update --no-tty
