# roles/appflowy/templates/appflowy-redis.container.j2
[Container]
Image=docker.io/library/redis:7-alpine
Exec=redis-server --requirepass {{ appflowy_redis_password }}
Network=appflowy.network

[Service]
Restart=always

[Install]
WantedBy=default.target
