# roles/systemd/templates/atuin-daemon.service.j2
[Unit]
Description=Atuin shell history sync daemon
After=network.target

[Service]
ExecStart=atuin daemon start --force
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
