# roles/systemd/templates/dotfiles-update.timer.j2
[Unit]
Description=Daily dotfiles update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
