# roles/systemd/templates/dotfiles-update.service.j2
[Unit]
Description=Dotfiles auto-update
{% if is_laptop | bool %}
After=network-online.target graphical-session.target
Wants=network-online.target
PartOf=graphical-session.target
{% else %}
After=network-online.target
Wants=network-online.target
{% endif %}

[Service]
Type=oneshot
# Wait for network to settle and desktop to stabilise before pulling
ExecStartPre=/bin/sleep {{ dotfiles_update_delay_seconds | default(60) }}
# Pull latest repo changes, then run ansible locally (no ansible-pull clone overhead).
# Timer cannot unlock Bitwarden non-interactively, so we explicitly skip secrets.
# SSH key sync happens on manual `dots-apply` instead.
ExecStart={{ brew_prefix }}/bin/just -f {{ ansible_facts['env']['HOME'] }}/.local/share/dotfiles/Justfile apply-nosecrets
# Also update podman containers while we're at it
ExecStartPost=-/usr/bin/podman auto-update --user
ExecStartPost=-/usr/bin/sudo /usr/bin/podman auto-update
Environment="PATH={{ brew_prefix }}/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME={{ ansible_facts['env']['HOME'] }}"
# Keep output for journalctl -u dotfiles-update
StandardOutput=journal
StandardError=journal

[Install]
{% if is_laptop | bool %}
WantedBy=graphical-session.target
{% else %}
WantedBy=default.target
{% endif %}
