# roles/netbox/files/netbox-postgres.container
[Unit]
Description=NetBox PostgreSQL
After=network-online.target

[Container]
Image=docker.io/library/postgres:18-alpine
ContainerName=netbox-postgres
Network=netbox.network
Volume=netbox-postgres.volume:/var/lib/postgresql/data
Environment=POSTGRES_DB=netbox
Environment=POSTGRES_USER=netbox
EnvironmentFile=/etc/netbox/postgres.env
# Disable Unix sockets — NetBox connects via TCP, and /var/run/postgresql
# may not be writable inside the container
Exec=postgres -c unix_socket_directories=''
AutoUpdate=registry

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
