# roles/netbox_agent/templates/netbox-agent.service.j2
[Unit]
Description=NetBox Agent - sync hardware inventory to NetBox
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/podman run --rm --privileged --net=host --pid=host \
          -v /etc/netbox:/etc/netbox:ro \
          -v /usr/sbin:/usr/sbin:ro \
          -v /usr/bin:/usr/bin:ro \
          -v /dev/mem:/dev/mem:ro \
          docker.io/python:3.12-slim \
          sh -c "pip install --no-cache-dir netbox-agent && netbox_agent -c /etc/netbox/netbox_agent.yml -r"
