# roles/n8n/templates/n8n-db.container.j2
[Container]
ContainerName=n8n-db
Image=docker.io/library/postgres:16-alpine
Environment=POSTGRES_USER=n8n
Environment=POSTGRES_PASSWORD={{ n8n_db_password }}
Environment=POSTGRES_DB=n8n
Volume=n8n-db-data:/var/lib/postgresql/data:Z
Network=n8n.network

[Service]
Restart=always

[Install]
WantedBy=default.target
