# roles/monitoring/files/grafana.container
[Unit]
Description=Grafana — metrics and logs visualization
After=monitoring-network.service prometheus.service loki.service
Wants=monitoring-network.service prometheus.service loki.service

[Container]
Image=docker.io/grafana/grafana-oss:latest
ContainerName=grafana
Network=monitoring.network

PublishPort=3001:3000

Volume=grafana-data.volume:/var/lib/grafana:z
Volume=%h/dev/homelab-monitoring/grafana/provisioning:/etc/grafana/provisioning:ro,z

# Auth — open access (homelab, Tailscale-protected)
Environment=GF_AUTH_ANONYMOUS_ENABLED=true
Environment=GF_AUTH_ANONYMOUS_ORG_ROLE=Admin
Environment=GF_AUTH_DISABLE_LOGIN_FORM=true
Environment=GF_SECURITY_ALLOW_EMBEDDING=true

# Branding
Environment=GF_DEFAULT_INSTANCE_NAME=Homelab
Environment=GF_SERVER_ROOT_URL=https://bihar.manatee-basking.ts.net/grafana
Environment=GF_SERVER_SERVE_FROM_SUB_PATH=true

# Home dashboard (fleet overview)
Environment=GF_DASHBOARDS_DEFAULT_HOME_DASHBOARD_PATH=/etc/grafana/provisioning/dashboards/overview.json

# Disable telemetry / news feed
Environment=GF_ANALYTICS_REPORTING_ENABLED=false
Environment=GF_ANALYTICS_CHECK_FOR_UPDATES=false
Environment=GF_NEWS_FEED_ENABLED=false

# Unified alerting
Environment=GF_ALERTING_ENABLED=false
Environment=GF_UNIFIED_ALERTING_ENABLED=true

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
