# roles/monitoring/files/perses.container
[Unit]
Description=Perses — metrics visualization dashboard
After=monitoring-network.service prometheus.service
Wants=monitoring-network.service prometheus.service

[Container]
Image=docker.io/persesdev/perses:latest
ContainerName=perses
Network=monitoring.network

PublishPort=8080:8080

Volume=perses-data.volume:/perses:z

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
