# roles/bluefin_common/templates/bluefin-wallpaper-rotate.timer.j2
[Unit]
Description=Rotate Bluefin wallpaper monthly timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
